# supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/etc/init.d/mysql start
autostart=true
autorestart=true

[program:httpd]
command=/etc/init.d/apache2 start
autostart=true
autorestart=true

